# tests/systemd/bar.socket
[Unit]
Description=Bar Socket

[Socket]
ListenStream=/run/bar.sock

[Install]
WantedBy=sockets.target